# pwr-ctrl.service
[Unit]
Description=Server power control

[Service]
ExecStart=/usr/bin/pwr-ctrl
Type=simple

[Install]
WantedBy=multi-user.target
